# tle-backup-service/systemd/tle-backup.service
[Unit]
Description=TLE-gf user.db off-site backup (one-shot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/tle-backup/backup.env
ExecStart=/opt/tle-backup-service/.venv/bin/python /opt/tle-backup-service/backup_user_db.py
# Don't hang a run forever if the network/source stalls.
TimeoutStartSec=600
# Light hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# Private writable /tmp so tempfile/locking work even under ProtectSystem=strict.
PrivateTmp=true
# Allow writing only where backups + the pinned known_hosts live.
ReadWritePaths=/var/backups/tle-gf /opt/tle-backup-service
